# tree/stage2/02-net-tweaks/files/internet-check.service
[Unit]
Description=Internet connectivity checker
After=network.target

[Service]
ExecStart=/usr/local/bin/internet-check forever
StandardOutput=journal
