# tree/stage2/02-net-tweaks/files/internet-check.service
[Unit]
Description=Internet connectivity checker
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/internet-check
