# tree/stage2/02-net-tweaks/files/internet-check.timer
[Unit]
Description=Internet connectivity checker

[Timer]
OnBootSec=1min
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
