# tree/stage2/02-net-tweaks/files/iptables-restore.service
[Unit]
Description=Restore iptables configuration
PartOf=offspot-runtime.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/find /etc/iptables/ -name '*.rules' -exec /sbin/iptables-restore {} \;

[Install]
WantedBy=multi-user.target
